# deploy/templates/igdb-sync.timer.j2
[Unit]
Description=Run IGDB sync nightly

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
