# deploy/templates/steamspy-enrich.service.j2
[Unit]
Description=SteamSpy popularity enrichment
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/igdb-scrape
ExecStart=/root/igdb-scrape/node_modules/.bin/tsx src/steamspy-enrich.ts
EnvironmentFile=/root/igdb-scrape/.env
StandardOutput=append:/var/log/igdb-scrape/steamspy-enrich.log
StandardError=append:/var/log/igdb-scrape/steamspy-enrich.log
