# deploy/templates/steamspy-enrich.timer.j2
[Unit]
Description=Run SteamSpy enrichment nightly

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
